# daemon/dbus-org.bootc.service
[Unit]
Description=bootc daemon

[Service]
BusName=org.bootc
ExecStart=/usr/libexec/bootc-daemon
